# lib/s6/init-run/service/.s6-svscan/finish
#!/bin/execlineb -S0

cd /
redirfd -w 2 /dev/console
fdmove -c 1 2
wait -r -- { }
/usr/lib/s6/init3 ${@}

